# filamentbox-webui.service
# Version: 1.6.0
[Unit]
Description=FilamentBox Web UI Server
After=network.target filamentbox.service
Requires=filamentbox.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/filamentcontrol
Environment="PATH=/opt/filamentcontrol/filamentcontrol/bin"
ExecStart=/opt/filamentcontrol/filamentcontrol/bin/python /opt/filamentcontrol/webui_server.py
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=filamentbox-webui

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/filamentcontrol

[Install]
WantedBy=multi-user.target
