# install/filamentbox.service.example
[Unit]
Description=FilamentBox Environment Monitor
Documentation=https://github.com/jdelgado-dtlabs/filamentenvmonitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=filamentbox
Group=filamentbox
WorkingDirectory=/opt/filamentcontrol

# ========================================
# Vault Configuration (Option 1 - Token Auth)
# ========================================
# Uncomment these lines if using Vault with token authentication:
# Environment="VAULT_ADDR=https://vault.example.com:8200"
# Environment="VAULT_TOKEN=your-vault-token"
# Optional for Vault Enterprise:
# Environment="VAULT_NAMESPACE=your-namespace"

# ========================================
# Vault Configuration (Option 2 - AppRole Auth - RECOMMENDED)
# ========================================
# Uncomment these lines if using Vault with AppRole authentication:
# Environment="VAULT_ADDR=https://vault.example.com:8200"
# Environment="VAULT_ROLE_ID=your-role-id"
# Environment="VAULT_SECRET_ID=your-secret-id"
# Optional for Vault Enterprise:
# Environment="VAULT_NAMESPACE=your-namespace"

# ========================================
# Alternative: Environment File
# ========================================
# Instead of inline Environment variables, you can use an environment file:
# EnvironmentFile=/etc/filamentbox/vault.env
# 
# Create /etc/filamentbox/vault.env with:
#   VAULT_ADDR=https://vault.example.com:8200
#   VAULT_ROLE_ID=your-role-id
#   VAULT_SECRET_ID=your-secret-id
# 
# Then set permissions:
#   sudo chmod 600 /etc/filamentbox/vault.env
#   sudo chown filamentbox:filamentbox /etc/filamentbox/vault.env

# ========================================
# Main Service Configuration
# ========================================
ExecStart=/opt/filamentcontrol/filamentbox/bin/python /opt/filamentcontrol/filamentbox.py

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=0

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/filamentcontrol

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=filamentbox

[Install]
WantedBy=multi-user.target
